# etc/udev/rules.d/99-PGenerator_bluez.rules
# Bluetooth Configuration
ACTION=="add", KERNEL=="hci[0-9]*", RUN+="/usr/bin/PGenerator_bluez.sh add %k"
